# deploy/systemd/backend.service
[Unit]
Description=Telecom SCM Backend
After=network.target mysql.service

[Service]
Type=simple
User=APP_RUN_USER
Group=APP_RUN_GROUP

WorkingDirectory=/opt/app/backend/current
ExecStart=/usr/bin/java -jar /opt/app/backend/current/telecom-scm-backend.jar

# 生产环境变量（也可改为读取 /opt/app/backend/.env 文件，由 systemd 加载）
Environment="SPRING_PROFILES_ACTIVE=db"
Environment="SERVER_ADDRESS=0.0.0.0"
Environment="SERVER_PORT=8080"
Environment="FILE_STORAGE_DIR=/opt/app/uploads"
# 请根据实际数据库修改
Environment="DB_URL=jdbc:mysql://127.0.0.1:3306/telecom_supply_chain?useUnicode=true&characterEncoding=UTF-8&serverTimezone=Asia/Shanghai&allowPublicKeyRetrieval=true&useSSL=false"
Environment="DB_USERNAME=root"
Environment="DB_PASSWORD=CHANGE_ME"
Environment="JWT_SECRET=CHANGE_ME_TO_A_LONG_RANDOM_STRING"

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
